# loom-intel-bridge.service
[Unit]
Description=LOOM weft-intel-bridge MCP Agent
After=network.target weft-intel.service
Wants=weft-intel.service

[Service]
Type=simple
User=gostev
WorkingDirectory=/mnt/workspace/projects/loom
ExecStart=/mnt/workspace/projects/loom/venv/bin/loom run agents/weft-intel-bridge.yaml --port 8200
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
