# loom-dashboard.service
[Unit]
Description=LOOM Web Dashboard
After=network.target loom-intel-bridge.service
Wants=loom-intel-bridge.service

[Service]
Type=simple
User=gostev
WorkingDirectory=/mnt/workspace/projects/loom
ExecStart=/mnt/workspace/projects/loom/venv/bin/python loom_dashboard.py 8300
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
